# dz_qinzi_baby.conf
[program:dz_qinzi_baby]
command=python3 dz_qinzi_baby.py
directory=/root/shangjia/dazhongdianping
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"